# deploy/com.curran.owlapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.curran.owlapi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/currangehring/Dev/owl-api/.venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>owl_api.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8787</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/currangehring/Dev/owl-api</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/currangehring/Library/Logs/owl-api/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/currangehring/Library/Logs/owl-api/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/currangehring/Dev/owl-api/.venv/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
